# gofn.service
[Unit]
Description=gofn snapshot

[Service]
Type=oneshot
ExecStart=/bin/sh /root/snapshot.service/snapshot.sh.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target